# systemd/gpsd.service
[Unit]
Description=GPS (Global Positioning System) Daemon

[Service]
ExecStart=/usr/sbin/gpsd -N

[Install]
Also=gpsd.socket
